# recipes-retroarch/retroarch/retroarch/init
#!/sbin/openrc-run

description="Start RetroArch"

depend()
{
    after coreui pulseaudio
    need localmount
    keyword -retroarch -libretro
}

start_pre()
{
    [ -f "/etc/profile.d/sysenv.sh" ] && source /etc/profile.d/sysenv.sh
    [ -f "/etc/profile.d/pulse.sh" ] && source /etc/profile.d/pulse.sh
}

command="retroarch"
command_args="-v -c /etc/retroarch.cfg --fullscreen --log-file=/var/log/retroarch.log"
command_background=true
pidfile="/run/retroarch.pid"
supervisor=supervise-daemon
